# sls/consul/files/consul.service
# Managed by Salt
[Unit]
ConditionFileNotEmpty=/etc/consul.d/main-config.json
